# deploy/kvllm.service.in
# kvllm — systemd *user* service template.
# Rendered to ~/.config/systemd/user/kvllm.service by deploy/install.sh (@UV@ / @WORKDIR@
# are filled in at install time). Runs as your user; with linger enabled it starts at boot.
# Manage via the `just service-*` recipes or `systemctl --user ... kvllm`.

[Unit]
Description=kvllm — vLLM OpenAI-compatible server (kai RTX 5090)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
WorkingDirectory=@WORKDIR@
# Model key, port, GPU fraction, and any HF token live here (see deploy/kvllm.env.example).
EnvironmentFile=@WORKDIR@/deploy/kvllm.env
# Resolve the registry key into `vllm serve ...` and exec into it (execvp ⇒ SIGTERM hits vLLM).
ExecStart=@UV@ run python -m kvllm.registry serve ${KVLLM_MODEL_KEY}
Restart=on-failure
RestartSec=5
# vLLM cold start = weights load + torch.compile/CUDA-graph capture; don't call it a failed start.
TimeoutStartSec=900
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=default.target
