# deploy/kvllm-helper.service.in
# kvllm-helper — systemd *user* service for the web control panel.
# Rendered to ~/.config/systemd/user/kvllm-helper.service by deploy/install.sh.
# Runs as your user (so it can drive `systemctl --user`), binds 0.0.0.0 for LAN access.

[Unit]
Description=kvllm helper — web control panel (model registry + switch + status)
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
WorkingDirectory=@WORKDIR@
EnvironmentFile=@WORKDIR@/deploy/kvllm.env
ExecStart=@UV@ run --group helper uvicorn kvllm.helper:app --host 0.0.0.0 --port ${KVLLM_HELPER_PORT}
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=default.target
